# systemd/greenclip.service
[Unit]
Description=Greenclip Clipboard manager

[Service]
ExecStart=%h/.local/bin/greenclip daemon
